# cron/crontab
* * * * * /cron/collect_logs.sh >> /proc/1/fd/1 2>> /proc/1/fd/2

